# systemd/ice.service
[Unit]
Description=rclone SFTP mount for Cooper ICE server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/home/jaeho/jaehho/.env
ExecStart=/home/jaeho/jaehho/scripts/ice/mount.sh
ExecStop=/bin/fusermount3 -u /mnt/ice
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
