# systemd/ice.service.tmpl
[Unit]
Description=rclone SFTP mount for Cooper ICE server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%REPO_ROOT%/.env
ExecStart=%REPO_ROOT%/scripts/ice/mount.sh
ExecStop=/bin/fusermount3 -u /mnt/ice
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
